# root/etc/services.d/caddy/run
#!/usr/bin/with-contenv sh

set -ex

exec s6-setuidgid abc caddy run -config /config/Caddyfile
